# app/supervisord.conf
[supervisord]
nodaemon=true                   ; Run supervisord in the foreground
logfile=/dev/null               ; Log supervisord output to stdout/stderr
logfile_maxbytes=0
loglevel=info

[program:redis]
command=/usr/bin/redis-server --bind 127.0.0.1 --port 6379 --loglevel notice
user=appuser                    ; Run redis as our non-root user
autorestart=true
priority=10
startsecs=3                     ; Wait 3 seconds to consider started
stdout_logfile=/dev/stdout      ; Redirect redis stdout to container stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr      ; Redirect redis stderr to container stderr
stderr_logfile_maxbytes=0

[program:gunicorn]
command=/bin/bash start.sh
directory=/app                  ; Working directory for the app
user=appuser                    ; Run gunicorn as our non-root user
autorestart=true
priority=20
startsecs=15                    ; Wait longer for gunicorn to start
depends_on=redis                ; Wait for redis to be running
environment=PYTHONUNBUFFERED=1,REDIS_HOST=127.0.0.1,REDIS_PORT=6379  ; Ensure Python output and Redis connection
stdout_logfile=/dev/stdout      ; Redirect gunicorn stdout to container stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr      ; Redirect gunicorn stderr to container stderr
stderr_logfile_maxbytes=0

# Optional: Add filebeat or other logging agents here if needed